# deploy/systemd/product-check.service
[Unit]
Description=ProductCheck Flask API (uWSGI)
After=network.target mysqld.service
Wants=mysqld.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/topuser/productCheck/backend
Environment=FLASK_ENV=production
Environment=PYTHONUTF8=1
Environment=PYTHONIOENCODING=utf-8
# EnvironmentFile=/home/topuser/productCheck/backend/.env

# 使用系统已安装的 uWSGI（Python 3.12）
ExecStart=/usr/local/python312/bin/uwsgi \
    --ini /home/topuser/productCheck/deploy/uwsgi/product_check.ini

ExecReload=/bin/kill -s HUP $MAINPID
KillSignal=SIGQUIT
TimeoutStopSec=30
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
